# downloadsort.service
[Unit]
Description=Sort Download Files

[Service]
ExecStart=/usr/bin/DownloadSort.sh

[Install]
WantedBy=multi-user.target
